# trusty/neutron-openvswitch/templates/git/upstart/neutron-plugin-openvswitch-agent.upstart
description "{{ service_description }}"
author "Juju {{ charm_name }} Charm <juju@localhost>"

start on runlevel [2345] and started {{ cleanup_process_name}}
stop on runlevel [!2345]

respawn

chdir /var/run

pre-start script
  mkdir -p /var/run/neutron
  chown neutron:root /var/run/neutron
end script

exec start-stop-daemon --start --chuid neutron --exec {{ executable_name }} -- \
    --config-file=/etc/neutron/neutron.conf --config-file={{ plugin_config }} \
    --log-file={{ log_file }}
